# externals/svscan.service
[Unit]
Description=svscan service
After=local-fs.target network.target

[Service]
ExecStart=/usr/bin/svscanboot
ExecStop=/usr/bin/pkill -f svscanboot
Type=simple
NonBlocking=yes

[Install]
WantedBy=multi-user.target
